# overlays/etc/systemd/system/activate-update.path
[Unit]
Description=Watch for application update requests

[Path]
# Trigger when this file is created
PathExists=/home/pi/.update_requested

[Install]
WantedBy=multi-user.target
