# overlays/etc/systemd/system/amaru-hotspot.service
[Unit]
Description=Reconcile Amaru fallback hotspot
Wants=NetworkManager.service
After=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/amaru-hotspot.sh reconcile